# misc/systemd/piwik-archive.timer
[Timer]
OnStartupSec=20m
OnUnitInactiveSec=30m

