# misc/systemd/piwik-archive.timer
[Install]
WantedBy=timers.target

[Timer]
OnStartupSec=20m
OnUnitInactiveSec=30m

